# roles/dns/files/etc/systemd/system/dnscrypt-proxy.socket.d/override.conf
[Socket]
ListenStream=
ListenDatagram=
ListenStream=[::1]:5353
ListenDatagram=[::1]:5353

