# roles/backup_server/files/etc/systemd/system/syncthing@backup.service.d/50-backup.conf
[Service]
Environment=HOME=/
ExecStart=
ExecStart=/usr/bin/syncthing -no-browser -no-restart -logflags=0 -home=/var/local/backup/syncthing
