# roles/dns/files/etc/systemd/system/dnscrypt-proxy.service.d/override.conf
[Service]
User=dns
ExecStart=
ExecStart=/usr/bin/dnscrypt-proxy /etc/dnscrypt/dnscrypt-proxy.conf

